# dotfiles/etc/profile.d/git-config-global-from-env.sh
# shellcheck shell=sh
. /usr/local/share/ubuntu-development/git-config-global-from-env.sh
